# ngrok.service
[Unit]
Description=ngrok tunnel service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi

# Start ngrok tunnel for port 8001
ExecStart=/usr/bin/env ngrok http 8001 --log stdout

# Restart on failure
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ngrok

# Environment - use pi user's home for config
Environment="HOME=/home/pi"

[Install]
WantedBy=multi-user.target
